# templates/lib/systemd/system/dsc-statistics-collector.service.erb
[Unit]
Description=A DNS Statistics Collector - Collector component

[Service]
Type=simple
ExecStart=/usr/bin/dsc -f -p /etc/dsc-statistics/dsc-collector.cfg

[Install]
WantedBy=multi-user.target
